# data/sp-daemon.socket
[Socket]
ListenStream=/tmp/spawny.socket

[Install]
WantedBy=sockets.target